# linchpin/on-demand-queue/on-demand-1.service
[Unit]
Description=Kickstart tests queue on-demand-1

[Service]
Type=simple
ExecStart=/home/kstest/on-demand-1/queue.sh
