# linchpin/on-demand-queue/on-demand-1-status.service
[Unit]
Description=Update status of on-demand-1 queue

[Service]
Type=simple
ExecStart=/home/kstest/on-demand-1/queue-update-status.sh
